# pfe_ddos_mitigation/webapp/ddos_mitigation.service
[Unit]
Description=Service de mitigation DDoS
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/suricata/pfe_ddos_mitigation/webapp
Environment="PATH=/home/suricata/pfe_ddos_mitigation/webapp/venv/bin"
ExecStart=/home/suricata/pfe_ddos_mitigation/webapp/venv/bin/gunicorn --workers 4 --bind 0.0.0.0:8000 wsgi:application
Restart=always

[Install]
WantedBy=multi-user.target
